# couchpotato/supervisor.ini
[supervisord]
nodaemon=true

[program:couchpotato]
user=couchpotato
directory=/CouchPotatoServer
priority=2
redirect_stderr=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
environment=HOME="/home/couchpotato",PWD="/CouchPotatoServer",LOGNAME="couchpotato",USER="couchpotato",TERM="xterm"
startsecs=5
command=/usr/bin/python /CouchPotatoServer/CouchPotato.py --config_file=/config/config.ini --data_dir=/config/data
