# deluge/1.3-stable/rootfs/etc/supervisor.d/supervisor.ini
[supervisord]
nodaemon=true

[program:deluged]
user=deluge
environment=TERM="xterm"
command=/usr/bin/deluged -d -c /config -L info -l /config/deluged.log

[program:deluge-web]
user=deluge
environment=PYTHON_EGG_CACHE="/home/deluge/.python-eggs",TERM="xterm"
command=/usr/bin/deluge-web -c /config -L info -l /config/deluge-web.log
