# samba/rootfs/etc/supervisor.d/supervisor.ini
[supervisord]
nodaemon=true

[program:nmbd]
command=/usr/sbin/nmbd -F

[program:smbd]
command=/usr/sbin/smbd -F

