# sickrage/supervisor.ini
[supervisord]
nodaemon=true

[program:sickrage]
user=sickrage
directory=/home/sickrage
priority=2
redirect_stderr=true
environment=HOME="/home/sickrage",PWD="/SickRage",LOGNAME="sickrage",USER="sickrage",TERM="xterm"
startsecs=5
command=/usr/bin/python /SickRage/SickBeard.py --datadir=/config
